# containers/zen/service/tracker/run
#!/bin/sh -e

PATH=/usr/local/bin:/usr/local/sbin:/bin:/sbin:/usr/bin:/usr/sbin:/usr/X11R6/bin
export HOME=/home/zen
cd /home/zen/nodetracker

exec 2>&1
exec node app.js >> log/zen-tracker.log
